# docker/db-backup/templates/crontab
0 0 * * * /backup.sh >> "/var/log/$(date +"%F").log"
